# dist/openrc/amdfan.in
#!/sbin/openrc-run

description="amdfan controller"
pidfile="/var/run/${RC_SVCNAME}.pid"

command="@bindir@/amdfan"
command_args="daemon --logfile=/var/log/amdfan.log --pidfile=${pidfile}"
command_args_background="--background"

extra_started_commands="reload"

reload() {
	start-stop-daemon --signal SIGHUP --pidfile "${pidfile}"
}

