# dist/openrc/amdfan.in
#!/sbin/openrc-run

description="amdfan controller"
command="@bindir@/amdfan"
command_args="--daemon"
command_background=true
pidfile="/var/run/${RC_SVCNAME}.pid"

reload() {
	start-stop-daemon --signal SIGHUP --pidfile "${pidfile}"
}

