# dist/systemd/amdfan.service
[Unit]
Description=amdfan controller
After=multi-user.target
Requires=multi-user.target

[Service]
ExecStart=/usr/bin/amdfan --daemon
Restart=always

[Install]
WantedBy=final.target

