# dist/systemd/amdfan.service.in
[Unit]
Description=amdfan controller
After=multi-user.target
Requires=multi-user.target

[Service]
ExecStart=@bindir@/amdfan --daemon
Restart=always

[Install]
WantedBy=final.target

