# dist/runit/log/run
#!/bin/sh
exec vlogger -t amdfan -p daemon
